# etc/systemd/system/mineos-agent.service
# /etc/systemd/system/mineos-agent.service
#
# Servizio principale: lancia il miner nativo (T-Rex/lolMiner/SRBMiner)
# puntato allo stratum Kryptex usando le credenziali in /opt/mineos/config.
# Lo script mineos-agent.sh legge rig.conf/wallet.conf/pools.conf e fa exec
# del miner corretto in foreground (così systemd ne gestisce il ciclo di vita).

[Unit]
Description=mineOS mining agent (miner nativo -> pool Kryptex)
Wants=network-online.target
After=network-online.target mineos-firstboot.service
# Non avviare prima che il first boot sia completato.
ConditionPathExists=/opt/mineos/state/first-boot.done
# Non avviare se serve un reboot (driver/kernel appena aggiornati).
ConditionPathExists=!/opt/mineos/state/reboot-required

[Service]
Type=simple
# Il miner gira in foreground; questo wrapper fa "exec" del binario scelto.
ExecStart=/opt/mineos/bin/mineos-agent.sh
# Riavvio robusto: se il miner crasha, riparte.
Restart=always
RestartSec=15
# Evita loop di restart troppo aggressivi.
StartLimitIntervalSec=300
StartLimitBurst=5
# Il mining è CPU/GPU intensivo: priorità I/O bassa per non affamare il sistema.
Nice=5
# Log su journald (consultabile con: journalctl -u mineos-agent -f).
StandardOutput=journal
StandardError=journal
# --- Hardening (compatibile con accesso GPU) ---
# NB: il mining richiede accesso ai device GPU, quindi NON usiamo PrivateDevices.
NoNewPrivileges=yes
ProtectSystem=strict
# Consenti scrittura solo dove serve (log/stato).
ReadWritePaths=/opt/mineos/logs /opt/mineos/state
ProtectHome=yes
ProtectControlGroups=yes
ProtectKernelModules=yes
RestartKillSignal=SIGINT

[Install]
WantedBy=multi-user.target
